# supporting_files/supervisord-actualizacion.conf
[program:actualizacion]
command=/actualizacion.sh
numprocs=1
autostart=true
autorestart=true