# supporting_files/supervisord-backup.conf
[program:backup]
command= bash /backup.sh
autostart=true
autorestart=true
priority=3