# supporting_files/supervisord-mysqld.conf
[program:mysqld]
command=/start-mysqld.sh
numprocs=1
priority=2
autostart=true
autorestart=true