# supporting_files/supervisord-actualizacion.conf
[program:actualizacion]
command= bash /actualizacion.sh
autostart=true
autorestart=true
priority=4