# supporting_files/supervisord-backup.conf
[program:backup]
command=/backup.sh
numprocs=1
priority=4
autostart=true
autorestart=true