# supporting_files/supervisord-actualizacion.conf
[program:actualizacion]
command=/actualizacion.sh